# assets/daemontools/log/run
#! /bin/sh
exec multilog t e 2>&1 ./main
